# configs/system/linta-weekly-snapshot.timer
# Linta Linux — Weekly Btrfs snapshot timer
# Profiles: [All]
# Spec: README.md §3.2 — Monday 00:00 UTC

[Unit]
Description=Weekly Btrfs snapshot (Monday 00:00 UTC)

[Timer]
OnCalendar=Mon *-*-* 00:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
